# etc/profile.d/sh.local
#Add any required envvar overrides to this file, it is sourced from /etc/profile
export GDK_DPI_SCALE=1.25
